# stolon/files/stolonkeeper.in
#!/bin/sh

#
# PROVIDE: stolonkeeper
# REQUIRE: stolonsentinel
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf to enable stolon-keeper:
#
#  stolonkeeper_enable="YES"
#
. /etc/rc.subr

name=stolonkeeper
rcvar=stolonkeeper_enable

su_cmd="/usr/bin/su"

start_cmd="stolonkeeper_start"
stolonkeeper_user=${stolonkeeper_user:-pgsql}
stolonkeeper_log=${stolonkeeper_log:-"/var/log/stolon-keeper.log"}


load_rc_config $name

: ${stolonkeeper_enable:="NO"}


command=%%PREFIX%%/bin/stolon-keeper

# configuration of kepeer include from file
. %%PREFIX%%/etc/stolon/stolon-keeper.conf


stolonkeeper_start()
{
	touch ${stolonkeeper_log} && chown ${stolonkeeper_user} ${stolonkeeper_log}
	${su_cmd} ${stolonkeeper_user} -c "exec ${command}  >${stolonkeeper_log} 2>&1" &
}

run_rc_command "$1"

